# deploy/taxi-bot.service
[Unit]
Description=Taxi Zhukovo Telegram Bot
After=network.target

[Service]
Type=simple
User=taxi
Group=taxi
WorkingDirectory=/opt/taxi-zhukovo
Environment="PATH=/opt/taxi-zhukovo/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/taxi-zhukovo/venv/bin/python /opt/taxi-zhukovo/run.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/taxi-zhukovo/logs/bot.log
StandardError=append:/opt/taxi-zhukovo/logs/bot_error.log

# Ограничения ресурсов (опционально)
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target

